# dots/x11/.config/systemd/user/wallpaper.timer
[Unit]
Description=Change wallpapers every hour

[Timer]
OnActiveSec=5sec
OnUnitActiveSec=59min

